# rootfs/etc/services.d/mopidy/run
#!/bin/sh
#fdmove -c 2 1
#s6-envdir env

echo "Starting mopidy server..."
mopidy
