# rootfs/etc/services.d/snapcast/run
#!/usr/bin/with-contenv sh
#fdmove -c 2 1
#s6-envdir env

echo "Starting snapcast server..."
snapserver --stream pipe:///tmp/snapfifo?name=default
