# rootfs/etc/services.d/fs-watcher/run
#!/bin/sh
#fdmove -c 2 1
#s6-envdir env

echo "Starting fs-watcher"
python3 /opt/fs-watcher/fs-watcher.py /streams $(pwd)/kill.sh
